# modprobed_db.timer
[Unit]
Description=Update modprobed-db database

[Timer]
OnUnitActiveSec=20m

[Install]
WantedBy=timers.target
